# supervisord.conf
[program:cored]
command=python /usr/local/sbin/core-daemon

[program:sshd]
command=/usr/sbin/sshd -D -d
