# supervisord.conf
[program:cored]
command=python /usr/bin/core-daemon

[program:sshd]
command=/usr/sbin/sshd -D -e
